# bender/files/crontab
0 9-20 * * MON-FRI bender notify-all
*/10 9-20 * * MON-FRI bender auto-processing